# provision/chef/cookbooks/packerd/templates/supervisord.conf.erb
[supervisord]
nodaemon=true

[program:docker]
command=/bin/bash -c "source /bin/docker_run.sh && exec $DOCKER daemon --debug"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:packerd]
command=/bin/bash -c "/opt/go/bin/packerd-server --tls-certificate /etc/ssl/private/temp.crt --tls-key /etc/ssl/private/temp.key --port 64154 --tls-port 64155 --host `/sbin/ifconfig eth0 | grep inet\ addr | awk -F: '{print $2}' | awk '{print $1}'`"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[include]
files = /etc/supervisor/conf.d/*.conf

